# pull_twitter_data/userinfo_cron_job_example.cron
MAILTO="adrian.benton@gmail.com"
PATH=/usr/bin:/bin:/usr/local/bin:/opt/anaconda2/bin
*/20 * * * *  /exp/abenton/twitter_brand/pull_twitter_data/pull_brand_data_userinfo.sh
0 0 1 * * /exp/abenton/twitter_brand/pull_twitter_data/pull_brand_data_pasttweets.sh