# pull_twitter_data/userinfo_cron_job_example.cron
MAILTO="belle.carrell@gmail.com"
0 */6 * * * sh /home/hltcoe/acarrell/PycharmProjects/twitter_brand/pull_twitter_data/pull_brand_data.sh userinfo >/dev/null 2>&1
0 0 1 * * sh /home/hltcoe/acarrell/PycharmProjects/twitter_brand/pull_twitter_data/pull_brand_data.sh pasttweets >/dev/null 2>&1